# launchd/com.nvinnikov.tg-claude.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nvinnikov.tg-claude</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/uv</string>
        <string>run</string>
        <string>python</string>
        <string>-m</string>
        <string>tgclaude.bot</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/nvinnikov/Documents/Flybo/tg-claude</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/nvinnikov/Documents/Flybo/tg-claude/tg-claude.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/nvinnikov/Documents/Flybo/tg-claude/tg-claude.log</string>
</dict>
</plist>
